# dktox/etc/profile.d/docker.sh
export EDITOR=/usr/bin/vim
export SHELL=/bin/bash

