# dotfiles/etc/supervisor.d/gzleo.ini
[program:dastask]
command=python3 /home/gzleo/triline/python/dastask.py
user=gzleo
autorestart=false
directory=/tmp
redirect_stderr=true

[program:mitmis]
command=python2 /home/gzleo/triline/python/mitmis.py
user=gzleo
autorestart=false
directory=/tmp
redirect_stderr=true
